# crontab
## Start workflow
00 03 * * * nohup /home/omnilab/OmniDataHouse/workflow.sh
